# .config/systemd/user/setxkbmap.service
[Unit]
Description=Set XKB layout

[Service]
ExecStart=/usr/bin/setxkbmap us

[Install]
WantedBy=default.target
